# deploy/systemd/netprobe-alertd.service
[Unit]
Description=TestMachine Alerts sweep
After=network-online.target

[Service]
Type=oneshot
User=netprobe
Group=netprobe
ExecStart=/opt/netprobe/venv/bin/python /opt/netprobe/app/jobs/alertd.py
