# deploy/systemd/netprobe-flow-exporter@.service
[Unit]
Description=TestMachine Flow Exporter (softflowd) on %I
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=netprobe

# socket di controllo per diagnostica
ExecStartPre=/usr/bin/install -d -m 0755 -o root -g netprobe /run/softflowd

# NOTE:
# -D  -> resta in foreground (così systemd lo tiene "active")
# -t  -> timeout in forma nome=valore (qui: general=60, expint=15)
# -T full -> tracking completo
# -v 9 -> NetFlow v9 (OK per nfcapd)
ExecStart=/usr/sbin/softflowd -D -i %I -n 127.0.0.1:2055 -v 9 -T full \
          -t general=60 -t expint=15 \
          -c /run/softflowd/%I.ctl

Restart=on-failure
RestartSec=1s

# hardening blando che non rompe libpcap
ProtectSystem=full
ProtectHome=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_PACKET

[Install]
WantedBy=multi-user.target
