# deploy/systemd/netprobe-flow-exporter@.service.d/override.conf
[Unit]
After=network-online.target netprobe-flow-collector.service
Wants=network-online.target
Requires=netprobe-flow-collector.service

[Service]
Restart=on-failure
RestartSec=2
